# kafka/supervisor/kafka-manager.conf
[program:kafka-manager]
command=/usr/bin/start-kafka-manager.sh
autostart=true
autorestart=true
